# feeder/crontab
0 */6 * * * root uv run src/feeder.py > /proc/1/fd/1 2> /proc/1/fd/2
